# cryospark-supervisor.timer
[Unit]
Description=Delay start of cryosparc-superviser.service

[Timer]
OnBootSec=4min


[Install]
WantedBy=timers.target
