# service.timer
[Unit]
Description=Delay start of sample.service

[Timer]
OnBootSec=4min


[Install]
WantedBy=timers.target
